# modules/mesaides/files/ma-monitor.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid ubuntu

chdir /opt/mes-aides

exec node monitor-server.js
