# modules/mesaides/files/openfisca.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

setuid main
setgid main

chdir /home/main/mes-aides-ui/openfisca

exec /home/main/venv/bin/gunicorn api --config config.py
